# extra/build/debian_pkg/python3-wasp-launcher.service
[Unit]
Description=WASP-Launcher Service
After=network.target

[Service]
EnvironmentFile=-/etc/default/python3-wasp-launcher
ExecStart=/usr/sbin/wasp-launcher
User=wasp_daemon
Group=wasp_daemon
RuntimeDirectory=wasp-launcher

[Install]
WantedBy=mutli-user.target
Alias=wasp-launcher.service

